# util/config/init.d/utwired
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
  need dbus
}

start() {
  ifconfig eth0 11.0.1.5 netmask 255.255.255.0
  route add -net default gw 11.0.1.1
  ifconfig eth0 up
}

stop() {
  ifconfig eth0 down
}

restart() {
  stop
  sleep 2
  start
}
