# util/config/init.d/utwireless
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
  need dbus
}

start() {
  ifconfig wlan0 up
  wpa_supplicant -iwlan0 -Dwext -c/etc/wpa_supplicant/wpa_supplicant.conf -B
  sleep 1
  ifconfig wlan0 192.168.16.105 netmask 255.255.0.0
  route add -net default gw "192.168.0.1"
}

stop() {
  killall wpa_supplicant
  ifconfig wlan0 down
}

restart() {
  stop
  sleep 2
  start
}
